# macOS/tcim/com.local.restart-tcim.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.local.restart-tcim</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>__HOME__/bin/restart-tcim.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>43200</integer>
	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/restart-tcim.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/restart-tcim.log</string>
</dict>
</plist>
